# supervisor.conf
[supervisorctl]
serverurl = unix:///opt/netuitive-agent/supervisor.sock

[unix_http_server]
file=/opt/netuitive-agent/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
http_port = /opt/netuitive-agent/supervisor.sock
minfds = 1024
minprocs = 200
loglevel = info
logfile = /dev/stdout
logfile_maxbytes = 0
nodaemon = True
pidfile = /opt/netuitive-agent/supervisord.pid
logfile_backups = 0


[program:netuitive-agent]
command=/opt/netuitive-agent/embedded/bin/python /opt/netuitive-agent/bin/netuitive-agent --foreground --configfile /opt/netuitive-agent/conf/netuitive-agent.conf
stderr_logfile=/dev/stderr
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
priority=999
startsecs=5
startretries=6
environment=SENSORS_LIB=/opt/netuitive-agent/embedded/lib/libsensors.so

[program:netuitive-statsd]
command=/opt/netuitive-agent/embedded/bin/netuitive-statsd --foreground --configfile /opt/netuitive-agent/conf/netuitive-agent.conf start
stderr_logfile=/dev/stderr
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
startsecs=5
startretries=3
priority=998
exitcodes=0

[group:netuitive-agent]
programs=netuitive-agent,netuitive-statsd
